# debian/postinst
#!/bin/bash
set -e

CUEMS_VENV="/usr/lib/cuems"
ENGINE_SERVICES="cuems-controller-engine cuems-node-engine"

case "$1" in
    configure)
        # Verify virtual environment exists
        if [ ! -f "$CUEMS_VENV/bin/python3" ]; then
            echo "ERROR: Virtual environment not found at $CUEMS_VENV" >&2
            echo "Please ensure cuems-utils package is properly installed." >&2
            exit 1
        fi

        # Verify console scripts were installed
        for script in controller-engine node-engine; do
            if [ ! -f "$CUEMS_VENV/bin/$script" ]; then
                echo "WARNING: Console script $script not found at $CUEMS_VENV/bin/" >&2
            fi
        done

        # Reload systemd (service files from cuems-common)
        if [ -d /run/systemd/system ]; then
            systemctl daemon-reload || true
        fi

        echo "cuems-engine installed successfully."
        echo "Package: cuemsengine installed to $CUEMS_VENV/lib/python*/site-packages/"
        echo "Console scripts: $CUEMS_VENV/bin/controller-engine, $CUEMS_VENV/bin/node-engine, $CUEMS_VENV/bin/system-ports"
        echo "Service files provided by cuems-common package."

        # Enable and (re)start the services
        if [ -d /run/systemd/system ]; then
            for svc in $ENGINE_SERVICES; do
                if ! systemctl is-enabled "$svc" >/dev/null 2>&1; then
                    systemctl enable "$svc" || echo "WARNING: Failed to enable $svc service" >&2
                fi

                # $2 is the previously-configured version (set on upgrade, empty on fresh install)
                if [ -n "$2" ]; then
                    echo "Upgrade from $2 detected, restarting $svc..."
                    systemctl restart "$svc" || echo "WARNING: Failed to restart $svc service" >&2
                else
                    systemctl start "$svc" || echo "WARNING: Failed to start $svc service" >&2
                fi
            done
        fi
        ;;

    abort-upgrade|abort-remove|abort-deconfigure)
        ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
        ;;
esac

#DEBHELPER#

exit 0

